# bot.conf
[program:telegrambot]
direcory=/home/admins/bot
command=python3 app.py

autostart=true
autorestart=true
stderr_logfile=/home/admins/logfile_err.log
stdout_logfile=/home/admins/logifle.log

